# right-click-fix/etc/xdg/autostart/right-click-fix.desktop
[Desktop Entry]
Exec=/usr/bin/right-click-fix
Name=Right Click Fix
Terminal=false
Type=Application
